# docker/nginx/nginx-symfony-php-7.4/files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /etc/php/php-fpm.conf
user=root
numprocs=1

;; We want to redirect nginx logging, which is sent to access/error logs
;; that redirect to /dev/stdout and /dev/stderr (respectively) to supervisor
stdout_logfile=/dev/stdout
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
redirect_stderr=true